# deploy/box/onebrain-metadata-drop.service
[Unit]
Description=OneBrain metadata-endpoint egress DROP (persist across reboots; before Docker)
# G1-6: the cloud-init runcmd inserts these iptables rules once, in memory; they vanish on
# the first reboot, re-exposing the cloud metadata endpoint (which serves the full user-data
# — for the MC/operator box, its ENTIRE baked /opt/onebrain/.env). This oneshot re-applies
# BOTH on every boot. Ordered BEFORE docker.service so the container-path DOCKER-USER drop is
# in place BEFORE dockerd starts any restart:always container: an After=docker unit races the
# async container start, so a container could reach 169.254.169.254 in the window between
# dockerd creating the (empty) DOCKER-USER chain and the oneshot re-inserting the drop. On a
# pre-docker boot the chain does not exist yet, so ExecStart CREATES it (idempotent) and inserts
# the drop; dockerd then finds DOCKER-USER already populated and only adds the FORWARD jump — it
# never flushes existing user rules, so the drop wins the race deterministically.
Before=docker.service
Wants=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
# -w waits for the xtables lock (nft-safe) so a dockerd race can't fail the drop.
ExecStart=/bin/sh -c 'iptables -w -C OUTPUT -d 169.254.169.254 -j DROP 2>/dev/null || iptables -w -I OUTPUT -d 169.254.169.254 -j DROP'
# DOCKER-USER may not exist yet on a pre-docker boot; CREATE it (idempotent) so the drop is in
# place before dockerd starts any container. dockerd preserves existing DOCKER-USER user rules
# on start (it only ensures the chain + the FORWARD jump), so this -I rule survives.
ExecStart=/bin/sh -c 'iptables -N DOCKER-USER 2>/dev/null || true; iptables -w -C DOCKER-USER -d 169.254.169.254 -j DROP 2>/dev/null || iptables -w -I DOCKER-USER -d 169.254.169.254 -j DROP'

[Install]
WantedBy=multi-user.target
