# deploy/box/onebrain-mc-image-prune.service
[Unit]
Description=Mission Control frequent rollback-safe Docker image reclaim
After=docker.service network-online.target onebrain-update.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
# Mission Control is deployed by hand (docker compose --force-recreate against a
# hand-edited images.override.yml, because MC only accepts production-signed
# releases on the auto path), so update.sh's immediate post-update prune never runs
# for it and the daily 168h onebrain-host-maintenance sweep is far too lax for its
# deploy cadence -- stale api+admin-ui+workers image trios accumulate until the
# root disk fills and the next --force-recreate wedges on "no space left on device".
#
# Reuse the SAME vetted, rollback-safe prune the daily timer runs (it protects every
# digest named by the compose file, both image overrides, the last applied release,
# and every running/stopped container; `docker image rm` is never forced), copied
# fresh from the running, digest-pinned Core API image so it always matches the
# deployed release. This MC-only unit differs from the daily one ONLY in cadence
# (its .timer) and retention (the MAINTENANCE_*_MIN_AGE_HOURS below).
#
# Copy to a private /run path so it never races the daily unit's /opt/onebrain copy,
# and point the data-volume verifier at its real /opt/onebrain location explicitly:
# the script otherwise resolves the verifier as a sibling of its own path, which
# would be a nonexistent /run/onebrain-data-volume.sh and make it hold (skip cleanup).
ExecStartPre=/usr/bin/docker compose --project-name {{COMPOSE_PROJECT}} -f /opt/onebrain/docker-compose.yml cp onebrain-api:/app/deploy/box/onebrain-host-maintenance.sh /run/onebrain-mc-image-prune.sh
ExecStartPre=/usr/bin/chmod 0755 /run/onebrain-mc-image-prune.sh
ExecStart=/usr/bin/env MAINTENANCE_IMAGE_MIN_AGE_HOURS=48 MAINTENANCE_BUILD_CACHE_MIN_AGE_HOURS=48 ONEBRAIN_DATA_VOLUME_VERIFY_SCRIPT=/opt/onebrain/onebrain-data-volume.sh /run/onebrain-mc-image-prune.sh
ExecStopPost=-/usr/bin/rm -f /run/onebrain-mc-image-prune.sh
Nice=19
IOSchedulingClass=idle
