# deploy/box/onebrain-drive-backup.timer
[Unit]
Description=Create a daily encrypted OneBrain database and Drive backup

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
